# deploy/symphony.service
[Unit]
Description=Symphony AI Agent Orchestrator
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=symphony
Group=symphony
WorkingDirectory=/opt/symphony

# Path to WORKFLOW.md — adjust to your deployment
ExecStart=/opt/symphony/.venv/bin/symphony run /opt/symphony/WORKFLOW.md \
    --logs-root /var/log/symphony \
    --log-level INFO

# Graceful shutdown: send SIGTERM and wait up to 60s before SIGKILL
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Restart policy — always restart unless manually stopped
Restart=on-failure
RestartSec=10

# Environment — secrets via environment variables
# Set LINEAR_API_KEY and GITHUB_TOKEN here or in a credentials file
EnvironmentFile=-/etc/symphony/env

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=symphony

[Install]
WantedBy=multi-user.target
